# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true

[program:shellinaboxd]
command=/usr/bin/shellinaboxd -t -p 80
autostart=true
